# mcboatface/conf/api_endpoint_supervisor.conf
[program:api_endpoint]
command=python /code/api/api.py 0.0.0.0 8000